# applications/luci-app-v2ray-server/root/etc/init.d/v2ray_server
#!/bin/sh /etc/rc.common
# Copyright (C) 2018-2020 Lienol <lawlienol@gmail.com>

START=99

CONFIG=v2ray_server
CONFIG_PATH=/var/etc/$CONFIG
LOG_PATH=/var/log/$CONFIG
LOG_APP_FILE=$LOG_PATH/app.log

echolog() {
	echo -e "$(date "+%Y-%m-%d %H:%M:%S"): $1" >> $LOG_APP_FILE
}

gen_v2ray_config_file() {
	config_get enable $1 enable
	[ "$enable" = "0" ] && return 0
	config_get remarks $1 remarks
	config_get port $1 port
	config_get transport $1 transport
	lua /usr/lib/lua/luci/model/cbi/v2ray_server/api/genv2rayconfig.lua $1 > $CONFIG_PATH/$1.json
	echolog "$remarks $port 生成并运行 V2ray 配置文件 - $CONFIG_PATH/$1.json" 
	if [ ! -f /var/v2server ]; then
		local ret="/usr/bin/v2ray"
		[ ! -f "$ret" ] && ret="/usr/bin/xray"
		cp -a $ret /var/v2server
	fi
	/var/v2server -config $CONFIG_PATH/$1.json >/dev/null 2>&1 &
}

start_v2ray_server() {
	mkdir -p $CONFIG_PATH $LOG_PATH
	touch $LOG_APP_FILE
	config_foreach gen_v2ray_config_file "user"
	fw3 reload > /dev/null 2>&1
}

stop_v2ray_server() {
	fw3 reload > /dev/null 2>&1
	ps -w | grep "$CONFIG_PATH/" | grep -v "grep" | awk '{print $1}' | xargs kill -9 >/dev/null 2>&1 &
	rm -rf $CONFIG_PATH
	rm -rf $LOG_PATH
}

status_v2ray_server() {
    config_load $CONFIG
    local global_enable=$(uci get $CONFIG.@global[0].enable 2>/dev/null)
    [ -z "$global_enable" ] && global_enable=0
    local running_processes=$(ps -w | grep "$CONFIG_PATH/" | grep -v "grep" | wc -l)
    local config_exists=0
    [ -d "$CONFIG_PATH" ] && config_exists=1
    echo "V2ray Server 状态信息："
    echo "-------------------------"
    echo "全局启用状态: $( [ "$global_enable" = "1" ] && echo "启用" || echo "禁用" )"
    echo "运行中实例数: $running_processes"
    echo "配置文件目录: $( [ $config_exists -eq 1 ] && echo "存在 ($CONFIG_PATH)" || echo "不存在" )"
    echo "日志文件路径: $LOG_APP_FILE"
    if [ $running_processes -gt 0 ]; then
        echo -e "\n运行中的进程详情:"
        ps -w | grep "$CONFIG_PATH/" | grep -v "grep" | awk '{print "PID: " $1 ", 配置文件: " $NF}'
    fi
}

start() {
	config_load $CONFIG
	enable=$(uci get $CONFIG.@global[0].enable)
	if [ "$enable" = "0" ];then
		stop_v2ray_server
	else
		start_v2ray_server
	fi
}

stop() {
	stop_v2ray_server
}

restart() {
	stop
	sleep 1
	start
}

status() {
    status_v2ray_server
}
